# systemd/user/getgovroam/getgovroam-notifs.timer
[Unit]
Description=Run getgovroam-notifs.service daily

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
